# scheduler/macos/com.earningstranscripts.job-scheduler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.earningstranscripts.job-scheduler</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/arunath/src/learn/earningstranscripts/.venv/bin/job-scheduler</string>
        <string>start</string>
        <string>--foreground</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/arunath/src/learn/earningstranscripts</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/arunath/src/learn/earningstranscripts/.venv/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>EARNINGS_DATA_DIR</key>
        <string>/Users/arunath/src/learn/data</string>
        <key>SCHEDULER_CONFIG_PATH</key>
        <string>/Users/arunath/src/learn/earningstranscripts/scheduler/examples/advanced_schedule.json</string>
    </dict>

    <!-- Start on login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart if it crashes -->
    <key>KeepAlive</key>
    <true/>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/Users/arunath/src/learn/data/logs/scheduler-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/arunath/src/learn/data/logs/scheduler-stderr.log</string>
</dict>
</plist>
